# hongs-wms/env/crontabs

1 0 * * * /opt/HongsWMS/bin/app system.serve --debug 2 >> /dev/null 2>&1
0 * * * * /opt/HongsWMS/bin/app system 60min --debug 2 >> /dev/null 2>&1
* * * * * /opt/HongsWMS/bin/app system 01min --debug 2 >> /dev/null 2>&1
