# hongs-wms/env/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/etc/init.d/cron start
exitcodes=0,2,70
user=root

[program:hongs-wms]
command=/opt/HongsWMS/bin/start.sh
user=nobody
