# hongs-web/env/crontabs

1 0 * * * /opt/HongsWeb/bin/app system.serve --debug 2 >> /dev/null 2>&1
